# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:post-rebuild]
command=/bin/bash -c 'test -x /home/node/.openclaw/lea/workspace/bin/post-rebuild-init.sh && /home/node/.openclaw/lea/workspace/bin/post-rebuild-init.sh || echo "[post-rebuild] init script not found, skipping"'
user=node
autostart=true
autorestart=false
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/node"

[program:gateway]
command=/bin/bash -c 'if [ -n "$OPENCLAW_GATEWAY_PASSWORD" ]; then openclaw gateway --bind "$OPENCLAW_GATEWAY_BIND" --port "$OPENCLAW_GATEWAY_PORT" --auth password --password "$OPENCLAW_GATEWAY_PASSWORD"; else openclaw gateway --bind "$OPENCLAW_GATEWAY_BIND" --port "$OPENCLAW_GATEWAY_PORT" --token "$OPENCLAW_GATEWAY_TOKEN"; fi'
user=node
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/node",NODE_ENV="production"

[program:mcporter-daemon]
command=/bin/bash -c 'sleep 10 && mcporter daemon start --foreground --log'
user=node
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/node"
